# docker/vscode-automation/supervisord.conf
[supervisord]
nodaemon=true
logfile=/logs/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Log rotation
[program:logrotate]
command=/usr/sbin/logrotate -f /automation/config/logrotate.conf
autostart=false
autorestart=false
startsecs=0
